# docker-custom/root/etc/services.d/node-red/run
#!/usr/bin/with-contenv bash

echo "Starting Node-Red"
cd /usr/src/node-red
exec s6-setuidgid node-red npm start --cache /data/.npm -- --userDir /data